# debian/prerm
#! /bin/bash

# prerm script fot kosmos-headless

set -e

systemctl stop detect-display.service
systemctl disable detect-display.service

